# saignes-weather-current.service
[Unit]
Description=Saignes-en-Padaine hourly airport (METAR) current-conditions refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=josue
WorkingDirectory=/home/josue/saignes_en_padaine
ExecStart=/usr/bin/python3 /home/josue/saignes_en_padaine/weather_mqtt.py --current-only
StandardOutput=journal
StandardError=journal
