# saignes-ack-listener.service
[Unit]
Description=Saignes-en-Padaine MQTT pump-ack listener (always-on, independent of the forecast cycle)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=josue
WorkingDirectory=/home/josue/saignes_en_padaine
ExecStart=/usr/bin/python3 /home/josue/saignes_en_padaine/weather_mqtt.py --ack-listener
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
