# saignes-dashboard.service
[Unit]
Description=Saignes-en-Padaine irrigation dashboard
After=network.target

[Service]
Type=simple
User=josue
WorkingDirectory=/home/josue/saignes_en_padaine
ExecStart=/usr/bin/python3 /home/josue/saignes_en_padaine/manage.py runserver --noreload 0.0.0.0:8080
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
