# saignes-weather.service
[Unit]
Description=Saignes-en-Padaine weather fetch + MQTT publish (one-shot, triggered by saignes-weather.timer)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=josue
WorkingDirectory=/home/josue/saignes_en_padaine
ExecStart=/usr/bin/python3 /home/josue/saignes_en_padaine/weather_mqtt.py
StandardOutput=journal
StandardError=journal
